# etc/periodic/hourly/freshclam
#!/bin/sh
/usr/bin/freshclam --quiet --log=/var/log/clamav/freshclam-hourly.log
